# debian/postinst
#!/bin/sh
set -e

PACKAGE_NAME="dxrt-driver-dkms"
DEBIAN_CHANGELOG="/usr/share/doc/${PACKAGE_NAME}/changelog.Debian.gz"
DKMS_MODULE_VERSION=$(dpkg-parsechangelog -l"$DEBIAN_CHANGELOG" \
                                           --show-field Version \
                                           | head -n 1 \
                                           | sed 's/[[:space:]]//g') 


echo "DEBUG: Final DKMS_MODULE_VERSION is: '${DKMS_MODULE_VERSION}'" >&2

if [ -z "${DKMS_MODULE_VERSION}" ]; then
    echo "Warning: Could not determine DKMS_MODULE_VERSION from changelog. Using default '0.0.0'." >&2
    DKMS_MODULE_VERSION="0.0.0" # Fallback to a default version
fi

# The postinst script is invoked by dpkg. The $1 variable represents the type of package management event currently being executed.
# =========================================================================
# Perform actions based on the package installation type.
# $1 variable indicates the type of package installation/upgrade/removal.
# 'configure' is the main stage after installation or upgrade.
# =========================================================================

case "$1" in
    configure)
        echo "Configuring ${PACKAGE_NAME} DKMS module version ${DKMS_MODULE_VERSION}..."

        # 1. Attempt to unload existing modules (safety measure)
        # Check if modules are actually loaded before attempting removal
        # This prevents duplicate removal attempts during package upgrades (prerm already handles removal)
        echo "Checking and unloading existing modules if needed..."
        
        # Check if dxrt systemd service is enabled and stop it if active
        DXRT_SERVICE_WAS_ENABLED=false
        if systemctl is-enabled --quiet dxrt 2>/dev/null; then
            DXRT_SERVICE_WAS_ENABLED=true
            # Stop the service if it's currently active
            if systemctl is-active --quiet dxrt 2>/dev/null; then
                echo "Stopping dxrt systemd service..."
                systemctl stop dxrt || true
            fi
        fi
        
        # Only remove modules if they are actually loaded (avoid duplicate removal)
        if lsmod | grep -q "^dxrt_driver"; then
            echo "Removing loaded dxrt_driver module..."
            rmmod -f dxrt_driver || true
        fi
        
        if lsmod | grep -q "^dx_dma"; then
            echo "Removing loaded dx_dma module..."
            rmmod -f dx_dma || true
        fi

        # 2. Check for conflicting manual installation
        # If manual installation is detected, fail with clear error message
        KERNEL_VERSION="$(uname -r)"
        MODULES_DIR="/lib/modules/${KERNEL_VERSION}"
        MANUAL_INSTALL_DETECTED=false
        
        echo "Checking for conflicting manual installation..."
        for module in dxrt_driver.ko dx_dma.ko; do
            if find "$MODULES_DIR/extra" -type f -name "$module" 2>/dev/null | grep -q .; then
                MANUAL_INSTALL_DETECTED=true
                echo "ERROR: Found manually installed module: $module" >&2
                find "$MODULES_DIR/extra" -type f -name "$module" 2>/dev/null
            fi
        done
        
        if [ "$MANUAL_INSTALL_DETECTED" = true ]; then
            echo "" >&2
            echo "===================================================================" >&2
            echo "ERROR: Manual installation detected in /lib/modules/" >&2
            echo "===================================================================" >&2
            echo "" >&2
            echo "This will conflict with DKMS package installation." >&2
            echo "" >&2
            echo "Please uninstall manually installed modules first:" >&2
            echo "  cd /path/to/dx_rt_npu_linux_driver/modules" >&2
            echo "  sudo ./build.sh -c uninstall" >&2
            echo "" >&2
            echo "Then install via DKMS package:" >&2
            echo "  sudo ./build.sh -c install-package" >&2
            echo "" >&2
            exit 1
        fi
        echo "No conflicting manual installation found."

        # 3. Add the module to DKMS
        # This registers the source code located at /usr/src/${PACKAGE_NAME}-${DKMS_MODULE_VERSION} with DKMS.
        echo "Adding ${PACKAGE_NAME}/${DKMS_MODULE_VERSION} to DKMS..."
        if ! dkms add "${PACKAGE_NAME}/${DKMS_MODULE_VERSION}"; then
            echo "Error: 'dkms add' failed for ${PACKAGE_NAME}/${DKMS_MODULE_VERSION}." >&2
            exit 1
        fi

        # 4. Build the DKMS module
        # This compiles the module against the currently running kernel.
        echo "Building ${PACKAGE_NAME}/${DKMS_MODULE_VERSION}..."
        if ! dkms build "${PACKAGE_NAME}/${DKMS_MODULE_VERSION}"; then
            echo "Error: 'dkms build' failed for ${PACKAGE_NAME}/${DKMS_MODULE_VERSION}." >&2
            exit 1
        fi

        # 5. Install the DKMS module
        # This copies the compiled module to /lib/modules/<kernel_version>/updates/dkms/.
        # Set environment variable to prevent DKMS hooks from running during package update
        echo "Installing ${PACKAGE_NAME}/${DKMS_MODULE_VERSION}..."
        if ! DKMS_PACKAGE_UPDATE=1 dkms install "${PACKAGE_NAME}/${DKMS_MODULE_VERSION}"; then
            echo "Error: 'dkms install' failed for ${PACKAGE_NAME}/${DKMS_MODULE_VERSION}." >&2
            exit 1
        fi

        # 6. Create Modprobe configuration file
        MODPROBE_CONF="/etc/modprobe.d/dx_dma.conf"
        echo "Creating modprobe configuration for module load order..."
        echo "softdep dx_dma pre: dxrt_driver post:" | sudo tee "$MODPROBE_CONF" > /dev/null

        # 7. Update initramfs
        echo "Updating initramfs..."
        update-initramfs -u 
        
        # 9. Ensure virt-dma dependency is satisfied
        # Check if CONFIG_DMA_VIRTUAL_CHANNELS is available (builtin or module)
        echo "Checking virt-dma dependency..."
        KERNEL_CONFIG="/boot/config-$(uname -r)"
        VDMA_BUILTIN=false
        VDMA_MODULE=false
        
        if [ -f "$KERNEL_CONFIG" ]; then
            if grep -Eq '^CONFIG_DMA_VIRTUAL_CHANNELS=y' "$KERNEL_CONFIG"; then
                VDMA_BUILTIN=true
                echo "-> CONFIG_DMA_VIRTUAL_CHANNELS=y (built-in)"
            elif grep -Eq '^CONFIG_DMA_VIRTUAL_CHANNELS=m' "$KERNEL_CONFIG"; then
                VDMA_MODULE=true
                echo "-> CONFIG_DMA_VIRTUAL_CHANNELS=m (module)"
            fi
        fi
        
        # If it's a module, ensure it's loaded before loading our drivers
        if [ "$VDMA_MODULE" = "true" ]; then
            if ! lsmod | grep -Eq '^virt_dma\b'; then
                echo "-> Loading virt_dma module..."
                modprobe virt-dma 2>/dev/null || modprobe virt_dma 2>/dev/null || {
                    echo "Warning: Could not load virt_dma module. Driver may fail to load." >&2
                }
            else
                echo "-> virt_dma module already loaded"
            fi
        fi
        
        # 10. Load the modules
        # Load the installed driver modules into the kernel.
        echo "Loading modules..."
        modprobe dx_dma || {
            echo "Warning: Failed to load dx_dma module. Please check dmesg for errors." >&2
        }
        modprobe dxrt_driver || {
            echo "Warning: Failed to load dxrt_driver module. Please check dmesg for errors." >&2
        }
        
        # 11. Restart dxrt systemd service if it was enabled before
        if [ "$DXRT_SERVICE_WAS_ENABLED" = "true" ]; then
            echo "Restarting dxrt systemd service..."
            systemctl restart dxrt || {
                echo "Warning: Failed to Restart dxrt systemd service. Please check systemctl status dxrt." >&2
            }
        fi
        
    ;;

    # Handle other package management events.
    # dkms remove operations are typically handled in the prerm script.
    abort-upgrade|abort-remove|abort-deconfigure)
        # No specific actions needed in postinst for these cases.
    ;;

    *)
        # Handle unknown arguments for robustness.
        echo "postinst: Unknown argument '$1'" >&2
        exit 1
    ;;
esac

exit 0
